# overlay-zork/chromeos-base/chromeos-bsp-zork/files/90-xhci-quirks.rules
ACTION=="offline", SUBSYSTEM=="usb", ENV{ERROR}=="DEAD", RUN+="/lib/udev/xhci-restart.sh"
