# baseboard-whalebook/chromeos-base/whaleos-angl-init/files/99-whaleos-angl.rules
ACTION=="add", ATTR{name}=="angl", RUN+="/lib/udev/whaleos-angl-init.sh %k"
